# systemd/cloudflare-dns.timer
[Unit]
Description=This is the timer to set the schedule for dns-update

[Timer]
OnCalendar=*:0/5
Persistent=true

[Install]
WantedBy=timers.target

