# systemd/cloudflare-dns-notify.service
[Unit]
Description=Notify current home IP address

[Service]
Type=oneshot
ExecStart=/opt/cloudflare-dns/dns-update -c /etc/dns-update.conf
